# port/pkg-plist
bin/vmstatedctl
etc/rc.d/vmstated
@sample %%ETCDIR%%/testvm/bhyve.conf.sample
@sample %%ETCDIR%%/testvm/config.sample
@sample %%ETCDIR%%/testvm/start_network.sample
@sample %%ETCDIR%%/testvm/stop_network.sample
sbin/vmstated
share/man/man1/vmstatedctl.1.gz
share/man/man8/vmstated.8.gz
